# ips/jenkins.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<!-- 

Copyright (C) 2007 Paul Oswald

Permission is hereby granted, free of charge, to any person obtaining a copy
of this software and associated documentation files (the "Software"), to deal
in the Software without restriction, including without limitation the rights
to use, copy, modify, merge, publish, distribute, sublicense, and/or sell
copies of the Software, and to permit persons to whom the Software is
furnished to do so, subject to the following conditions:

The above copyright notice and this permission notice shall be included in
all copies or substantial portions of the Software.

THE SOFTWARE IS PROVIDED "AS IS", WITHOUT WARRANTY OF ANY KIND, EXPRESS OR
IMPLIED, INCLUDING BUT NOT LIMITED TO THE WARRANTIES OF MERCHANTABILITY,
FITNESS FOR A PARTICULAR PURPOSE AND NONINFRINGEMENT. IN NO EVENT SHALL THE
AUTHORS OR COPYRIGHT HOLDERS BE LIABLE FOR ANY CLAIM, DAMAGES OR OTHER
LIABILITY, WHETHER IN AN ACTION OF CONTRACT, TORT OR OTHERWISE, ARISING FROM,
OUT OF OR IN CONNECTION WITH THE SOFTWARE OR THE USE OR OTHER DEALINGS IN
THE SOFTWARE.
-->

<service_bundle type="manifest" name="Jenkins">
	<service name="application/jenkins" type="service" version="1">
		
		<create_default_instance enabled="true" />
		
		<!-- Only one instance of Jenkins should ever run per server -->
		<single_instance />
		
		<dependency name="multi-user-server" type="service" grouping="require_all" restart_on="none">
			<service_fmri value="svc:/milestone/multi-user-server" />
		</dependency>

		<method_context>
			<method_credential user='root' group='root' />
			<method_environment>
				<envvar name='PATH' value='/usr/bin:/usr/sbin:/usr/ccs/bin:/usr/local/bin:/usr/local/sbin:/usr/sfw/bin' />
				<envvar name='JENKINS_HOME' value='/var/lib/jenkins' />
			</method_environment>
		</method_context>

		<!-- Set the JENKINS_HOME env variable, and run the war file in /apps/jenkins/jenkins.war -->
		<exec_method type="method" name="start"	exec="java -Xmx512m -jar /usr/local/bin/jenkins.war" timeout_seconds="0"/>
		<exec_method type="method" name="stop" exec=":kill -TERM" timeout_seconds="30"/>

		<!-- We are going to be kicking off a single child process so we want Wait mode-->
		<property_group name='startd' type='framework'>
			<propval name='duration' type='astring' value='child' />
		</property_group>

		<stability value="Unstable" />

		<template>
			<common_name>
				<loctext xml:lang='C'>Jenkins Continuous Build Server</loctext>
			</common_name>
			<documentation>
              <doc_link name='jenkins-ci.org' uri='http://jenkins-ci.org/' />
			</documentation>
		</template>
	</service>
</service_bundle>
